# scripts/xarcade2jstick.service
[Unit]
Description=Xarcade2JstickESC

[Service]
Type=oneshot
ExecStart=/usr/local/bin/xarcade2jstickESC -s

[Install]
WantedBy=multi-user.target
